# VM/rc.local
#!/bin/bash

theta start --config=../privatenet/node --password="qwertyuiop"

exit 0
